# rpi-packages/board/rpi3/rootfs_overlay/var/spool/cron/crontabs/root
*/60 * * * * ntpdate -v -b in.pool.ntp.org
